# systemd/frey-chrome.service.example
[Unit]
Description=Chrome with CDP for Frey
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/google-chrome \
    --remote-debugging-port=9222 \
    --user-data-dir=/tmp/chrome-cdp-profile \
    --headless \
    --disable-gpu \
    --no-sandbox \
    --disable-dev-shm-usage
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
